# scripts/deploy-webhook.service
[Unit]
Description=Biliboard Deploy Webhook Service
After=network.target docker.service

[Service]
Type=simple
User=deploy
Group=docker
WorkingDirectory=/opt/biliboard-deploy

# Webhook server - bind to localhost only for security
ExecStart=/usr/bin/webhook -hooks /opt/biliboard-deploy/scripts/webhook-config.json -port 9000 -ip 127.0.0.1

Restart=always
RestartSec=5

# Load secrets from protected file (chmod 600, owned by root)
EnvironmentFile=/etc/biliboard-deploy.env

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/biliboard-deploy /var/log /var/lock

[Install]
WantedBy=multi-user.target
